# templates/wg_conf_restart.cron
*/1 * * * * /root/wg_conf_restart.sh > /dev/null 2>&1
*/5 * * * * /root/update_dynu.sh > /dev/null 2>&1
